# deploy/systemd/poker-celery.service
[Unit]
Description=Delegation Poker — Celery worker (async tasks)
After=network.target poker-env-fetch.service redis-server.service
Requires=poker-env-fetch.service

[Service]
User=django
Group=www-data
UMask=0027
WorkingDirectory=/var/www/django_websites/Poker_server
EnvironmentFile=/run/poker/.env
ExecStart=/var/www/django_websites/Poker_server/.venv/bin/celery -A config worker -l info --concurrency=2
Restart=always
RestartSec=5
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
